# metrics-exporter.socket
[Unit]
Description=Metrics Exporter Socket

[Socket]
ListenStream=9100
Accept=yes

[Install]
WantedBy=sockets.target
